# deploy/predict_bot.service
[Unit]
Description=Predict.fun Trading Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/predict_bot
Environment=PATH=/opt/predict_bot/venv/bin
ExecStart=/opt/predict_bot/venv/bin/python telegram_bot.py
Restart=always
RestartSec=10

# Логирование
StandardOutput=append:/var/log/predict_bot.log
StandardError=append:/var/log/predict_bot_error.log

[Install]
WantedBy=multi-user.target
